# src/upstart/swap.conf
start on started udevtrigger

console output

pre-start script
    echo "Activating all swap files/partitions..."
end script

exec /sbin/swapon -a
