# src/upstart/sysctl.conf
start on stopped mount-kernel-filesystems

console output

script
	echo "Starting sysctl...."
	/sbin/sysctl -q -p
	echo "sysctl DONE"
end script
