# src/upstart/tty2.conf
# tty2 - getty
#
# This service maintains a getty on tty2 from the point the system is
# started until it is shut down again.

start on started udevtrigger
stop on shutdown


respawn
exec /sbin/agetty 9600 tty2
