# src/upstart/welcome.conf
start on startup

console output

script
	echo "Welcome to IPFire $(cat /etc/system-release)"
end script
