# src/upstart/welcome.conf
description	"Set sysctls from /etc/sysctl.conf"
author		"Scott James Remnant"

start on startup

console output

script
	echo "Welcome to IPFire $(cat /etc/system-release)"
end script
