# src/systemd/firewall6.service.in
[Unit]
Description=Firewall for IPv6
After=firewall-init.service
Before=network.target
Requires=firewall-init.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/firewall6 start
ExecStop=@sbindir@/firewall6 stop

[Install]
WantedBy=multi-user.target
