# src/systemd/network-init.service.in
[Unit]
Description=Basic Initialization for Network Connectivity.
After=firewall-init.service
Before=network.target
Wants=network.target
Requires=firewall-init.service

[Service]
Type=oneshot
ExecStart=@sbindir@/network init

[Install]
WantedBy=multi-user.target
