# src/systemd/network@.service.in
[Unit]
Description=Network Connectivity for zone %I
After=firewall-init.service network-init.service
Requires=firewall-init.service network-init.service
Wants=network.target
Before=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/network zone %I up
ExecStop=@sbindir@/network zone %I down

[Install]
WantedBy=network.target
Alias=network.target.wants/network@%i.service
