# src/initscripts/extras/collectd.conf
description		"Start system statistics daemon"
author			"IPFire Team"

start on stopped mountfs
stop on starting shutdown

exec /usr/sbin/collectd -C /etc/collectd.conf
