# src/initscripts/extras/fcron.conf
description		"Start the cron daemon"
author			"IPFire Team"

start on stopped mountfs
stop on shutdown

exec /usr/sbin/fcron -y -b
expect fork
respawn
