# src/initscripts/extras/lighttpd.conf
description		"Start lighttpd"
author			"IPFire Team"

start on started network
stop on stopping network

exec /usr/sbin/lighttpd -f /etc/lighttpd.conf
expect daemon
