# src/initscripts/extras/openldap.conf
description		"Starts OpenLDAP slapd daemon"
author			"IPFire Team"

start on started mountfs
stop on shutdown or reboot

exec /usr/sbin/slapd
expect fork
respawn
