# src/initscripts/extras/pdns-recursor.conf
description		"Control the powerdns recursor"
author			"IPFire Team"

start on starting network
stop on stopped network

exec /usr/sbin/pdns_recursor
expect fork
respawn
