# src/initscripts/extras/rstp.conf
description		"Rapid Spanning Tree Protocol Daemon"
author			"IPFire Team"

start on starting network
stop on stopped network

pre-start script
	modprobe net_pf_17
end script

exec /sbin/rstpd
