# src/initscripts/extras/syslog-ng.conf
description		"Start system log daemon"
author			"IPFire Team"

start on stopped mountfs
stop on shutdown

exec /sbin/syslog-ng -F
respawn
