# src/initscripts/extras/vsftpd.conf
description		"vsftpd FTPS Server"
author			"IPFire Team"

start on started network
stop on stopping network

exec /usr/sbin/vsftpd
respawn
