# config/udev/60-net.rules
KERNEL=="eth*", NAME="port%n"
SUBSYSTEM=="net", RUN+="/etc/init.d/networking/net-hotplug"
