# config/udev/60-net.rules
SUBSYSTEM=="net", RUN+="/etc/init.d/networking/net-hotplug"
