# pkgs/core/initscripts/init/cleanfs.conf
description		"Cleaning / filesystem"
author			"IPFire Team"

start on stopping mountfs

script
	rm -rf /tmp/* /var/{lock,run}/*
end script
