# pkgs/core/initscripts/init/swap.conf
description		"Activating swap partitions"
author			"IPFire Team"

start on stopped checkfs

exec /sbin/swapon -a
