# pkgs/core/initscripts/init/tty1.conf
description		"Start a tty"
author			"IPFire Team"

start on stopped mountfs
stop on starting shutdown or starting reboot

respawn
exec /sbin/agetty 9600 tty1
