# src/install/etc/init/tty1.conf
# tty1 - getty
#
# This service maintains a getty on tty1 from the point the system is
# started until it is shut down again.

start on stopped mountfs
stop on shutdown

respawn
exec /sbin/agetty -n -l /root/autologin tty1 9600
