# src/install/etc/init/tty1.conf
description		"Start a tty"
author			"IPFire Team"

start on stopped mountfs
stop on shutdown

respawn
exec /sbin/agetty -n -l /root/autologin tty1 9600
