# src/install/etc/init/tty2.conf
description		"Start a tty"
author			"IPFire Team"

start on stopped mountfs
stop on starting shutdown

respawn
exec /sbin/agetty -n -l /root/autologin tty2 9600
