# src/install/etc/init/tty3.conf
description		"Start a tty"
author			"IPFire Team"

start on stopped mountfs
stop on shutdown

respawn
exec /sbin/agetty -n -l /root/autologin tty3 9600
