# src/initscripts/core/cleanfs.conf
description		"Cleaning / filesystem"
author			"IPFire Team"

start on stopping mountfs

pre-start script
	rm -rf /tmp/* /var/{lock,run}/*
end script
