# src/initscripts/core/control-alt-delete.conf
description	"Control-Alt-Delete / Handling"
author		"IPFire Team"

start on control-alt-delete

task
exec shutdown -r now "Control-Alt-Delete pressed"
