# src/initscripts/core/control-alt-delete.conf
description	"Control-Alt-Delete / Handling"
author		"IPFire Team"

start on control-alt-delete

exec initctl emit --no-wait reboot
