# src/initscripts/core/random
#!/bin/sh
# Begin /etc/init.d/random

### BEGIN INIT INFO
# Provides:            random
# Required-Start:      udev
# Should-Start:
# Required-Stop:       sendsignals
# Should-Stop:
# Default-Start:       S
# Default-Stop:        0 6
# Short-Description:   Saves and restores the random seed.
# Description:         Saves and restores the random seed.
# X-LFS-Default-Start: S25
# X-LFS-Default-Stop:  K45
# X-LFS-Provided-By:
### END INIT INFO

. /lib/lsb/init-functions

case "${1}" in
    start)
        message="Initializing kernel random number generator..."
        if [ -f /var/tmp/random-seed ]; then
            cat /var/tmp/random-seed > /dev/urandom
        fi
        dd if=/dev/urandom of=/var/tmp/random-seed count=1 &>/dev/null
        evaluate_retval standard
        ;;

    stop)
        message="Saving random seed..."
        dd if=/dev/urandom of=/var/tmp/random-seed count=1 &>/dev/null
        evaluate_retval standard
        ;;

    *)
        echo "Usage: ${0} {start|stop}"
        exit 1
        ;;
esac

# End /etc/init.d/random
