# src/initscripts/core/tty4.conf
description		"Start a tty"
author			"IPFire Team"

start on stopped mountfs
stop on starting shutdown

respawn
exec /sbin/agetty 9600 tty4
